# docs/solarpowered-main/files/solarpowered-ex/usr/lib/systemd/system/b550-suspend-fix.service
[Unit]
Description=B550 boards suspend fix

[Service]
ExecStart=/bin/sh -c "/bin/echo GPP0 > /proc/acpi/wakeup"

[Install]
WantedBy=multi-user.target
